# roles/trouter/templates/alfresco-transform-router.service
[Unit]
Description=Alfresco Transform Router
After=syslog.socket network.target local-fs.target remote-fs.target

[Service]
Type=simple
User={{ username }}
ExecStart={{ binaries_folder }}/ats-atr.sh
WorkingDirectory={{ ats_home }}
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
