# roles/transformation-services/templates/alfresco-transform.service
[Unit]
Description=Alfresco Transform Service
Requires=alfresco-tengine-aio.service
PartOf=alfresco-content.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=true

[Install]
WantedBy=multi-user.target